# scripts/launchd/com.cryptotrader.scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cryptotrader.scanner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/code/cryptotrader/scripts/launchd/run-scanner.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__HOME__/code/cryptotrader</string>

  <!-- Run every 30 minutes -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/cryptotrader/scanner.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/cryptotrader/scanner.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
