# scripts/launchd/com.cryptotrader.harvester-paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cryptotrader.harvester-paper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/elmo.asmussen/Projects/Crypto/scripts/launchd/run-harvester-paper.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/elmo.asmussen/Projects/Crypto</string>

  <!-- Daily at 02:15 local (after the 00:00 UTC daily-bar close). launchd
       runs a missed calendar fire when the machine next wakes. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <!-- Fire once at load too (idempotent — verifies the wiring). -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/elmo.asmussen/Library/Logs/cryptotrader-harvester/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/elmo.asmussen/Library/Logs/cryptotrader-harvester/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
